# initramfs/etc/inittab
::sysinit:/etc/rc
::respawn:/bin/sh
::ctrlaltdel:/sbin/reboot
